# supervisord.conf
; ================================
;  huey supervisor
; ================================
[supervisord]
user=nonroot
loglevel=info

[program:backup]
command=python pdfding/manage.py run_huey
stdout_logfile=background_tasks.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs=2

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true
